# docker-configs/worker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=999
stdout_events_enabled=true
stderr_events_enabled=true

[program:scheduler]
command=sh /var/www/scheduler.sh
numprocs=1
autostart=true
autorestart=true
user=root

[program:queue]
command=php /var/www/artisan queue:work redis --tries=1
numprocs=1
autostart=true
autorestart=true
user=root
stdout_logfile=/var/www/storage/logs/queue_out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/storage/logs/queue_err.log
stderr_logfile_maxbytes=0

[include]
files = /etc/supervisor/conf.d/*.conf
